# config/systemd/overrides/tollgate-csms.service.d/override.conf
# systemd drop-in override for tollgate-csms (OCPP 1.6 CSMS).
#
# The tollgate-csms binary (this repo) uses ocpp-go's cs.Start(port, path)
# which binds 0.0.0.0:port by default. Charge points reach the CSMS via
# caddy reverse-proxy on a public https URL, so the CSMS itself only needs
# to accept loopback connections from caddy.
#
# This drop-in uses systemd's cgroup BPF network filter to drop ALL
# non-loopback packets — external clients cannot reach port 8887 directly
# even though `ss` still shows the bind on *:8887. Charge points must
# connect via the caddy-fronted public URL.
#
# Install: copy this directory to /etc/systemd/system/tollgate-csms.service.d/
# Then: systemctl daemon-reload && systemctl restart tollgate-csms
[Service]
IPAddressDeny=any
IPAddressAllow=localhost
