# config/systemd/tollgate-settle.timer
[Unit]
Description=Tollgate Settlement — weekly Monday 03:00

[Timer]
OnCalendar=Mon *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
