# config/systemd/overrides/tollgate-webssh.service.d/override.conf
# systemd drop-in override for tollgate-webssh (WebSocket → SSH bridge).
#
# Forces the listen address to loopback. Caddy reverse-proxies public
# wss:// traffic to 127.0.0.1:8092 — the webssh binary never needs to
# accept external connections directly.
#
# Install: copy this directory to /etc/systemd/system/tollgate-webssh.service.d/
# Then: systemctl daemon-reload && systemctl restart tollgate-webssh
[Service]
Environment=TOLLGATE_WEBSSH_ADDR=127.0.0.1:8092
