# config/systemd/overrides/tollgate-net.service.d/override.conf
# systemd drop-in override for tollgate-net (session daemon).
#
# The tollgate-net binary (separate Rust repo) has no --bind flag, so it
# listens on 0.0.0.0:2121 by default. This drop-in uses systemd's cgroup
# BPF network filter to drop ALL non-loopback packets — external clients
# cannot reach the port even though `ss` still shows the bind on 0.0.0.0.
#
# Why this is safe: tollgate-net is only consumed by sibling services
# (tollgate-auth-radius via TOLLGATE_SESSIOND_URL=http://127.0.0.1:2121,
# tollgate-auth-ocpi via delegated mode). All callers are on loopback.
#
# Install: copy this directory to /etc/systemd/system/tollgate-net.service.d/
# Then: systemctl daemon-reload && systemctl restart tollgate-net
[Service]
IPAddressDeny=any
IPAddressAllow=localhost
